# browser-images/selenium-images/chromium/chrome/supervisord-vnc.conf
[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=2
loglevel=info
pidfile=/var/run/supervisor/supervisord.pid
minfds=1024
minprocs=200
childlogdir=/var/log/supervisor
user=jelenoid

[unix_http_server]
file=/var/run/supervisor/supervisor.sock
chmod=0700
chown=jelenoid:jelenoid

[supervisorctl]
serverurl=unix:///var/run/supervisor/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

# Xvfb — виртуальный framebuffer (display :99)
[program:xvfb]
command=/usr/bin/Xvfb :%(ENV_DISPLAY_NUM)s -screen 0 %(ENV_SE_SCREEN_WIDTH)sx%(ENV_SE_SCREEN_HEIGHT)sx%(ENV_SE_SCREEN_DEPTH)s -ac -nolisten tcp -dpi 96 +extension RANDR
autorestart=true
startretries=5
startsecs=1
stopwaitsecs=5
priority=10
stdout_logfile=/dev/null
stderr_logfile=/dev/null
user=jelenoid

# openbox — лёгкий WM, чтобы Chrome вёл себя как на реальном desktop
[program:openbox]
command=/usr/bin/openbox
environment=DISPLAY=":%(ENV_DISPLAY_NUM)s",HOME="/home/jelenoid",XAUTHORITY="/home/jelenoid/.Xauthority"
autorestart=true
startretries=5
startsecs=1
stopwaitsecs=5
priority=20
stdout_logfile=/dev/null
stderr_logfile=/dev/null
user=jelenoid

# x11vnc — публикация дисплея по VNC
[program:x11vnc]
command=/usr/bin/x11vnc -display :%(ENV_DISPLAY_NUM)s -forever -shared -rfbport %(ENV_SE_VNC_PORT)s -rfbauth /home/jelenoid/.vncpasswd
environment=DISPLAY=":%(ENV_DISPLAY_NUM)s",HOME="/home/jelenoid",XAUTHORITY="/home/jelenoid/.Xauthority"
autorestart=true
startretries=5
startsecs=1
stopwaitsecs=5
priority=30
stdout_logfile=/dev/null
stderr_logfile=/dev/null
user=jelenoid

[program:chromedriver]
command=/opt/chromedriver/chromedriver-linux64/chromedriver --port=4445 --url-base=/wd/hub --allowed-origins=*
environment=CHROME_BIN="/opt/chrome/chrome-linux64/chrome",DISPLAY=":99",XAUTHORITY="/home/jelenoid/.Xauthority"
directory=/home/jelenoid
autorestart=true
startretries=3
startsecs=1
stopwaitsecs=20
priority=40
stdout_logfile=/dev/null
stderr_logfile=/dev/null
user=jelenoid

# socat — прокси 0.0.0.0:4444 → 127.0.0.1:4445
[program:chromedriver-proxy]
command=/usr/bin/socat -d tcp-listen:4444,fork,reuseaddr tcp:127.0.0.1:4445
environment=DISPLAY=":99",XAUTHORITY="/home/jelenoid/.Xauthority"
autorestart=true
startretries=3
startsecs=1
stopwaitsecs=20
priority=50
stdout_logfile=/dev/null
stderr_logfile=/dev/null
user=jelenoid
